# supervisor/php_dev.conf
[program:nginx]
command=nginx

[program:mysql]
command=mysqld

[program:php71-fpm]
command=php-fpm7.1